# setup/ansible/roles/ai-server/templates/scheduled-reboot.timer.j2
[Unit]
Description=Scheduled reboot

[Timer]
OnCalendar=*-*-* {{ scheduled_reboot_time }}
Persistent=true

[Install]
WantedBy=timers.target