# setup/ansible/roles/ssh-tunnel-and-watchdog/templates/scheduled-reboot.service.j2
[Unit]
Description=Scheduled reboot

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reboot
